# AgBiz-Logic-dev/conf/systemd/dev/agbiz.service
[Unit]
Description=agbiz

[Service]
Slice=machine.slice
ExecStart=/home/hammonse/projects/AgBiz-Logic/app/rkt-run.sh
KillMode=mixed
Restart=always
